# debian/libravatar-slave.cron.d
*/15 * * * *     libravatar-slave       test -e /usr/bin/libravatar-update-mirror && /usr/bin/libravatar-update-mirror
