# debian/libravatar-www.cron.daily
#!/bin/sh
django-admin clearsessions --settings=libravatar.settings --pythonpath=/usr/share/libravatar

# Delete files that were uploaded more than a day ago
find /var/lib/libravatar/uploaded -type f -ctime +1 -delete

# Old exports need to be deleted after a day
find /var/lib/libravatar/export -type f -ctime +1 -delete

# Check that old files don't accumulate here
find /var/lib/libravatar/ready -type f -ctime +1 -print
