# debian/libravatar-master.cron.d
*/15 * * * *     libravatar-master       test -e /usr/bin/libravatar-update-mirrors && /usr/bin/libravatar-update-mirrors
