# debian/libravatar-www.cron.daily
#!/bin/sh
django-admin cleanup --settings=libravatar.settings --pythonpath=/usr/share/libravatar
